# deploy/passivediscovery.service
[Unit]
Description=Passive Network Discovery Service
Documentation=https://github.com/nvp/passivediscovery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=discovery
Group=discovery
WorkingDirectory=/opt/passivediscovery
ExecStart=/opt/passivediscovery/discovery \
    --interface ${DISCOVERY_INTERFACE} \
    --output ${DISCOVERY_OUTPUT} \
    --db ${DISCOVERY_DB} \
    --oui ${DISCOVERY_OUI} \
    --api-addr ${DISCOVERY_API_ADDR} \
    --ui \
    --log-level ${DISCOVERY_LOG_LEVEL} \
    --log-output ${DISCOVERY_LOG_OUTPUT}
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Logging — send to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=passivediscovery

# Security sandbox
NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/var/lib/passivediscovery
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true

# Environment (override via /etc/default/passivediscovery)
EnvironmentFile=-/etc/default/passivediscovery

[Install]
WantedBy=multi-user.target