# deploy/passivediscovery.service
[Unit]
Description=Passive Network Discovery Service
Documentation=https://github.com/nvp/passivediscovery
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=discovery
Group=discovery
WorkingDirectory=/opt/passivediscovery
ExecStart=/opt/passivediscovery/discovery \
    --interface ${DISCOVERY_INTERFACE} \
    --output ${DISCOVERY_OUTPUT} \
    --db ${DISCOVERY_DB} \
    --oui ${DISCOVERY_OUI} \
    --api-addr ${DISCOVERY_API_ADDR} \
    --ui \
    --log-level ${DISCOVERY_LOG_LEVEL} \
    --log-output ${DISCOVERY_LOG_OUTPUT}
Restart=on-failure
RestartSec=10

# Logging — send to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=passivediscovery

# Security sandbox
NoNewPrivileges=false
ProtectSystem=strict
ReadWritePaths=/var/lib/passivediscovery
PrivateTmp=true
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

# Defaults — these can be overridden in /etc/default/passivediscovery
Environment=DISCOVERY_INTERFACE=eth0
Environment=DISCOVERY_OUTPUT=/var/lib/passivediscovery/output
Environment=DISCOVERY_DB=/var/lib/passivediscovery/discovery.db
Environment=DISCOVERY_OUI=/opt/passivediscovery/oui.csv
Environment=DISCOVERY_API_ADDR=:8080
Environment=DISCOVERY_LOG_LEVEL=info
Environment=DISCOVERY_LOG_OUTPUT=stdout

# Override via /etc/default/passivediscovery (loaded after Environment=, takes precedence)
EnvironmentFile=-/etc/default/passivediscovery

[Install]
WantedBy=multi-user.target